# scripts/pi_systemd/vignette-capture.timer
[Unit]
Description=Trigger Vignette capture every 30 seconds while session is active

[Timer]
OnBootSec=15s
OnUnitInactiveSec=30s
AccuracySec=1s
Unit=vignette-capture.service

[Install]
WantedBy=timers.target
